# overlay/etc/xdg/autostart/calamares.desktop
[Desktop Entry]
Type=Application
Name=Calamares
Exec=pkexec calamares
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
