# overlay/etc/xdg/autostart/calamares.desktop
[Desktop Entry]
Type=Application
Name=Calamares
Exec=/usr/local/bin/start-installer
NoDisplay=true
